# gsisshd.service
[Unit]
Description=gsissh server daemon
Documentation=man:gsisshd(8) man:gsisshd_config(5)
After=syslog.target network.target auditd.service

[Service]
EnvironmentFile=/etc/sysconfig/gsisshd
ExecStartPre=/usr/sbin/gsisshd-keygen
ExecStart=/usr/sbin/gsisshd -D $OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target
